# Library/LaunchAgents/com.claude.memory-log.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.claude.memory-log</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/bash</string>
		<string>/Users/james.niu/claude-os/scripts/1-log.sh</string>
	</array>
	<key>StartInterval</key>
	<integer>3600</integer>
	<key>StandardOutPath</key>
	<string>/Users/james.niu/claude-os/output/1-log-stdout.log</string>
	<key>StandardErrorPath</key>
	<string>/Users/james.niu/claude-os/output/1-log-stderr.log</string>
	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>/Users/james.niu/.local/bin:/usr/local/bin:/usr/bin:/bin</string>
	</dict>
</dict>
</plist>
