# launchd/com.claude.memory-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.claude.memory-sync</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/bash</string>
		<string>$HOME/claude-os/scripts/4-sync-confluence.sh</string>
	</array>
	<key>StartInterval</key>
	<integer>86400</integer>
	<key>StandardOutPath</key>
	<string>$HOME/claude-os/output/4-sync-stdout.log</string>
	<key>StandardErrorPath</key>
	<string>$HOME/claude-os/output/4-sync-stderr.log</string>
	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>/usr/local/bin:/usr/bin:/bin</string>
		<key>CONFLUENCE_EMAIL</key>
		<string>james.niu@basis.com</string>
		<key>CONFLUENCE_TOKEN</key>
		<string>YOUR_TOKEN_HERE</string>
	</dict>
</dict>
</plist>
